# Python/varnish-nginx-gunicorn-flask-celery-postgres/application/conf/supervisor/celery.conf
[supervisord]
nodaemon=true

[program:celery_service]
directory=/application/flaskr/
startretries=10000
autostart=true
user=root
command=/bin/bash -c "celery -A app.celery worker -B"
